# .crontab
## Cubietruck Crontab File

*     *     *     *     *     crontab /hdd/.cubie/.crontab
*     *     *     *     *     find /hdd/.cubie -type f -exec dos2unix {} \;
