# home/cubie/.crontab
## Cubietruck Crontab File

# --------------------------------------------------
# System scripts
# --------------------------------------------------

*     *     *     *     *     crontab /mnt/ssd/.cubie/home/cubie/.crontab
*     *     *     *     *     chown -R cubie:cubie /mnt/ssd/.cubie
# *     *     *     *     *     chown -R root:root /mnt/ssd
# 0     0     *     *     *     reboot

# --------------------------------------------------
# Development scripts
# --------------------------------------------------

*     *     *     *     *     find /mnt/ssd/.cubie -type f -exec dos2unix {} \;
