# home/cubie/.crontab
## Cubietruck Crontab File

# --------------------------------------------------
# System scripts
# --------------------------------------------------

*     *     *     *     *     crontab /mnt/ssd/.cubie/home/cubie/.crontab
*     *     *     *     *     find /mnt/ssd/.cubie \! -user cubie -exec chown -R cubie:cubie {} \;
*     *     *     *     *     find /mnt/ssd/projects \! -user cubie -exec chown -R cubie:cubie {} \;
*     *     *     *     *     find /mnt/raid \! -user cubie -exec chown -R cubie:cubie {} \;

# --------------------------------------------------
# Development scripts
# --------------------------------------------------

*     *     *     *     *     find /mnt/ssd/.cubie -type f -exec dos2unix -ic {} \; | xargs dos2unix
*     *     *     *     *     find /mnt/ssd/projects -type f -exec dos2unix -ic {} \; | xargs dos2unix
0     10    *     *     *     bash /mnt/ssd/.cubie/etc/init.d/show_network_status
0     0     *     *     *     bash /mnt/ssd/.cubie/home/cubie/scripts/turn_off_leds
